# docker/crontab/crontab
* * * * * cd /var/www/codes && php artisan schedule:run >> /dev/stdout 2>&1
